# magento1-simple/devops/magento.cron
# Execute Magento cron engine every 5 minutes
*/5 * * * *     application  php /var/www/cron.php
